# etc/systemd/system/zapret-discord.service
[Unit]
Description=Zapret DPI Bypass for Discord Voice
After=network.target

[Service]
Type=simple
ExecStart=/opt/my-firewall-setup/opt/zapret/nfq/nfqws \
    --qnum=202 \
    --filter-udp=1-65535 \
    --hostlist=/opt/my-firewall-setup/opt/zapret/ipset/discord-hosts.txt \
    --dpi-desync=fake \
    --dpi-desync-fooling=badseq \
    --dpi-desync-repeats=5 \
    --dpi-desync-any-protocol=1
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
